# build/etc/service/rsyslog/run
#!/bin/bash
exec 1>&2

rm -f /run/rsyslogd.pid

exec /usr/sbin/rsyslogd -n
